# email_updates/email-updates.service
[Unit]
Description=Email available updates to root

[Service]
Type=oneshot

ExecStart=/usr/bin/email-updates.sh
